# packaging/deb/etc/systemd/system/netops-server.service
[Unit]
Description=NetOps Backend Server
Documentation=https://github.com/amishagrawal2001-arch/tlink-netops
After=network.target

[Service]
Type=simple
User=netops
Group=netops
WorkingDirectory=/opt/netops-server
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
Environment=PORT=4000
Environment=NODE_ENV=production
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
